# pollbot.service
[Unit]
Description=PollBot Discord Bot
After=network.target

[Service]
Type=simple
User=blink
WorkingDirectory=/home/blink/dev/pollbot
ExecStart=/home/blink/.nvm/versions/node/v24.11.1/bin/node dist/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
